# deploy/drinktea-backend.service
[Unit]
Description=DrinkTea Backend FastAPI Application
Documentation=https://github.com/yourusername/drinktea
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/drinktea/backend
Environment="PATH=/opt/drinktea/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# 启动命令 - 使用 Uvicorn ASGI 服务器
ExecStart=/opt/drinktea/backend/.venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-level info \
    --access-log \
    --no-use-colors

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 资源限制
LimitNOFILE=65536
TimeoutStartSec=30
TimeoutStopSec=30

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drinktea-backend

[Install]
WantedBy=multi-user.target
